# deploy/systemd/ultra-memory-embed.timer
[Unit]
Description=Embed memory anchors frequently (prevents multi-hour backlogs on Pi)

[Timer]
OnBootSec=1min
OnUnitActiveSec=2min
Persistent=true
AccuracySec=30s

[Install]
WantedBy=timers.target
